# app/trunk/i2v/main/crontabs/klogcln
#!/bin/sh

/bin/sed -i '/Unsupported command/ d' /var/log/messages
/bin/sed -i '/last message repeated 8 times/ d' /var/log/messages

